# receiver/kde-tags-announce.service
[Unit]
Description=kde-tags - mDNS announcement (_kdetags._tcp) on the local network

[Service]
ExecStart=@AVAHI_PUBLISH@ -s "@INSTANCE@" _kdetags._tcp 1 "name=@NAME@" "topic=@TOPIC@" "server=@SERVER@" "user=@USER@"
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target
